# data/export/custom_upstart/resque.conf.erb
# Resque script config
start on starting <%= app %>-<%= name %>
stop on stopping <%= app %>-<%= name %>

respawn

# Wait for job to finish processing then restart
reload signal QUIT

env PORT=<%= port %>
<% engine.env.each do |name,value| -%>
env <%= name.upcase %>='<%= value.gsub(/'/, "'\"'\"'") %>'
<% end -%>

<%
  project_path = engine.root.to_s.split("/")[0..-3].join("/")
  current_path = File.join(project_path, "current")
%>

setuid <%= user %>

chdir <%= project_path %>

script
  cd <%= current_path %>
  exec <%= process.command %>
end script
